# scripts/tests/fixtures/hermes-gateway-base.service
[Unit]
Description=Hermes gateway fixture

[Service]
Type=simple
ExecStart=/bin/true

[Install]
WantedBy=default.target
